# configuration/supervisord/golfstats.conf
[program:golfstats-uwsgi]
command=/opt/golfstats/venv/bin/uwsgi /opt/golfstats/app/configuration/uwsgi/golfstats.ini
autostart=true
autorestart=true
redirect_stderr=true